# docker/supervisord.conf
#######################################################################
# Supervisor config - DO NOT MODIFY UNLESS YOU KNOW WHAT YOU ARE DOING
#######################################################################
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/run/supervisord.pid
user=root

[unix_http_server]
file=/dev/shm/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock

#################################################################
# System Jobs - DO NOT MODIFY UNLESS YOU KNOW WHAT YOU ARE DOING
#################################################################
[program:go-grpc-websocket-test]
command=go-grpc-websocket-test
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:caddy]
command=caddy run --config /etc/Caddyfile
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0
